# src/etc/systemd/system/s7-updates-motd.service
[Unit]
Description=Update MOTD with available package updates

[Service]
ExecStart=/usr/local/bin/s7-updates-motd.sh
